# ansible/systemd/letsencrypt.timer
[Unit]
Description=Renew letsencrypt certificates

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target
