# ansible/systemd/letsencrypt.service
[Unit]
Description=Renew letsencrypt certificates

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
