# ansible/systemd/webinvoices-currency-rates.timer
[Unit]
Description=Import currency rates in webinvoices

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
